# php-fpm/supervisor.conf
[supervisord]
nodaemon = true
[program:myprocess1]
command=php-fpm
[program:myprocess2]
command=/usr/sbin/sshd -D 
[program:myprocess3]
command=systemctl start sshd